# crontab
# Fend Marketplace Cron Jobs
# This file should be placed in /etc/cron.d/fend-cron

# Run subscription expiry notifications every day at 9am
0 9 * * * root cd /app && /usr/local/bin/python manage.py send_subscription_expiry_notifications >> /var/log/cron.log 2>&1

# Run monthly deals digest on the 1st of every month at 10am
0 10 1 * * root cd /app && /usr/local/bin/python manage.py send_monthly_deals_digest >> /var/log/cron.log 2>&1

# Run payment verification alerts every 6 hours
0 */6 * * * root cd /app && /usr/local/bin/python manage.py send_payment_verification_alerts >> /var/log/cron.log 2>&1

# (This empty line is required by cron)
